# etc/systemd/system/github-runner.target
[Unit]
Description=the target for github-runner.service
Requires=github-runner.service systemd-logind.service systemd-user-sessions.service
