# etc/systemd/system/github-runner-config.service
[Unit]
Description=GitHub Actions Runner Config
After=network.target

[Service]
Type=oneshot
User=runner
EnvironmentFile=/etc/environment
EnvironmentFile=/etc/github-runner-env
WorkingDirectory=/home/runner
ExecStart=/home/runner/actions-runner-config.sh

[Install]
WantedBy=multi-user.target
